# nplayer.service
[Unit]
Description=Nativity Player

[Service]
Type=simple
ExecStart=/root/proj/player

[Install]
WantedBy=graphical.target
